# config/systemd/user/rclone-timer.service
[Unit]
Description=Rclone Periodic GDrive Bisync
Documentation=https://rclone.org/
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/bin/rclone_timer.sh
StandardOutput=append:%h/.config/rclone/rclone.log
StandardError=append:%h/.config/rclone/rclone.log
