# config/systemd/user/taskwarrior-notify.service
[Unit]
Description=Taskwarrior notification script

[Service]
Type=oneshot
ExecStart=%h/scripts/taskwarrior.sh
